# config/cron/production.crontab
# start unicorn on boot
@reboot  /home/eric/konklone.com/shared/cron/start.sh

# update sitemap each morning
0 8 * * * /home/eric/konklone.com/shared/cron/rake.sh sitemap

# daily db backups, monthly asset backups

@daily /home/eric/backup/mongo.sh > /home/eric/backup/output-mongo.txt 2>&1
@monthly /home/eric/backup/assets.sh > /home/eric/backup/output-assets.txt 2>&1

# keep trailing newline
